# docker/reddit/supervisord.conf
[supervisord]
nodaemon=true

# reddit web app
# use paster for development, gunicorn_paster for production
# TODO test gunicorn_paster
[program:gunicorn]
command=/home/reddit/src/reddit/scripts/wrap-job paster serve --reload /home/reddit/src/reddit/r2/run.ini
# command=/usr/bin/gunicorn_paster --log-file=/var/log/syslog /home/reddit/src/reddit/r2/run.ini

# reddit-media and reddit-pixel
[program:nignx]
command=/usr/sbin/nginx -g 'daemon off;'

# activity service
# redis could be own container, but have to fork https://github.com/reddit-archive/reddit-service-activity
[program:reddit-activity]
directory=/home/reddit/src/activity
command=/bin/bash -c "service redis-server start && export PYTHONPATH=$PYTHONPATH:/usr/local/lib/python2.7/site-packages:/usr/local/lib/python2.7/dist-packages && baseplate-serve2 --bind localhost:9002 /home/reddit/src/activity/example.ini"

# websocket service
# deprecated. not used and needs forked example.ini for rabbitmq connection.
# [program:reddit-websockets]
# directory=/home/reddit/src/websockets
# command=/bin/bash -c "export PYTHONPATH=$PYTHONPATH:/usr/local/lib/python2.7/site-packages:/usr/local/lib/python2.7/dist-packages && baseplate-serve2 --bind localhost:9001 /home/reddit/src/websockets/example.ini"

# click and geoip services
[program:reddit-gunicorn-services]
command=service gunicorn start

# rabbitmq consumers
# replaces manage-consumers and consumer-count.d consumer configuration. manually add more instances below.
[program:reddit-consumer-author_query_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle author_query_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/voting.py -c 'consume_author_query_queue()'"

[program:reddit-consumer-automoderator_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle automoderator_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/automoderator.py -c 'run()'"

[program:reddit-consumer-butler_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle butler_q1 $REDDIT_INI -c 'from r2.lib.butler import run; run()'"

[program:reddit-consumer-commentstree_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle commentstree_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/db/queries.py -c 'run_commentstree()'"

[program:reddit-consumer-del_account_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle del_account_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/db/queries.py -c 'consume_deleted_accounts()'"

[program:reddit-consumer-domain_query_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle domain_query_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/voting.py -c 'consume_domain_query_queue()'"

# TODO: consumer is broken, see issues/44
# [program:reddit-consumer-event_collector_q1]
# command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle event_collector_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/eventcollector.py -c 'from pylons import app_globals; process_events(app_globals, limit=100)'"

[program:reddit-consumer-markread_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle markread_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/db/queries.py -c 'consume_mark_all_read()'"

[program:reddit-consumer-newcomments_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle newcomments_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/db/queries.py -c 'run_new_comments()'"

[program:reddit-consumer-scraper_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle scraper_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/media.py -c 'run()'"

# TODO: is crashing, need solr container
# [program:reddit-consumer-search_q1]
# command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle search_q1 $REDDIT_INI -c 'from pylons import app_globals; app_globals.search.run_changed()'"

[program:reddit-consumer-subreddit_query_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle subreddit_query_q1 $REDDIT_INI $REDDIT_ROOT/r2/lib/voting.py -c 'consume_subreddit_query_queue()'"

[program:reddit-consumer-vote_comment_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle vote_comment_q1 $REDDIT_INI -c 'from r2.lib.voting import consume_comment_vote_queue; consume_comment_vote_queue()'"

[program:reddit-consumer-vote_link_q1]
command=/bin/bash -c "source /etc/default/reddit && exec /home/reddit/src/reddit/scripts/wrap-job paster run --proctitle vote_link_q1 $REDDIT_INI -c 'from r2.lib.voting import consume_link_vote_queue; consume_link_vote_queue()'"
