# docker/solr/supervisord.conf
[supervisord]
nodaemon=true

[program:solr]
command=/usr/share/tomcat7/bin/catalina.sh run
