# services/reddit-consumer-author_query_q@.service
[Unit]
Description=author_query_q instance %i